# infra/systemd/genergi-worker.service
[Unit]
Description=GENERGI Worker
After=network.target redis.service redis-server.service genergi-api.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/genergi/current
EnvironmentFile=/opt/genergi/shared.env
ExecStart=/usr/bin/env bash -lc 'cd /opt/genergi/current && corepack pnpm --filter @genergi/worker start'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
